# deploy/go2rtc.service.example
[Unit]
Description=go2rtc streaming gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/unifi-doorbell-monitor
ExecStart=/usr/local/bin/go2rtc -config /opt/unifi-doorbell-monitor/config/go2rtc.yaml
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
